# deploy/launchd/com.ccc.regress.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ccc.regress</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/SHARED/path/to/CCC/scripts/ccc-board.py</string>
    <string>regress</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/SHARED/path/to/business-app</string>
  <!-- MUST be an engine-eligible app path (not CCC orch). -->
  <!-- Install: copy to ~/Library/LaunchAgents/com.ccc.regress.plist and load. -->
  <!-- One plist per business app, or wrap a small shell that loops apps. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>23</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/ccc-regress.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/ccc-regress.err</string>
</dict>
</plist>
