# scripts/com.ccc.chat-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- TEMPLATE ONLY — do not load as-is.
     Generate machine-local paths via: bash scripts/install-hub-plist.sh [--start]
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ccc.chat-server</string>
  <key>ProgramArguments</key>
  <array>
    <string>CCC_HOME/.venv-hub/bin/python</string>
    <string>CCC_HOME/scripts/ccc-chat-server.py</string>
    <string>--no-open</string>
  </array>
  <key>WorkingDirectory</key>
  <string>CCC_HOME</string>
  <key>KeepAlive</key>
  <true/>
  <key>RunAtLoad</key>
  <true/>
  <key>ProcessType</key>
  <string>Background</string>
  <key>ThrottleInterval</key>
  <integer>10</integer>
  <key>EnvironmentVariables</key>
  <dict>
    <key>CCC_CHAT_HOST</key>
    <string>0.0.0.0</string>
    <key>CCC_CHAT_PORT</key>
    <string>7777</string>
    <key>CCC_CHAT_USER</key>
    <string>ccc</string>
    <key>CCC_CHAT_PASS</key>
    <string>ccc</string>
    <key>CCC_BOARD_URL</key>
    <string>http://127.0.0.1:7775</string>
    <key>CCC_CHAT_NO_OPEN</key>
    <string>1</string>
    <key>ANTHROPIC_BASE_URL</key>
    <string>https://api.minimaxi.com/anthropic</string>
    <key>ANTHROPIC_MODEL</key>
    <string>MiniMax-M3</string>
  </dict>
</dict>
</plist>
