# scripts/com.ccc.chat-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ccc.chat-server</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/python3</string>
    <string>/Users/apple/program/CCC/scripts/ccc-chat-server.py</string>
    <string>--no-open</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/apple/program/CCC</string>
  <key>KeepAlive</key>
  <true/>
  <key>RunAtLoad</key>
  <true/>
  <key>ProcessType</key>
  <string>Background</string>
  <key>ThrottleInterval</key>
  <integer>10</integer>
  <key>EnvironmentVariables</key>
  <dict>
    <key>CCC_CHAT_HOST</key>
    <string>0.0.0.0</string>
    <key>CCC_CHAT_PORT</key>
    <string>7777</string>
    <key>CCC_CHAT_USER</key>
    <string>ccc</string>
    <key>CCC_CHAT_PASS</key>
    <string>ccc</string>
    <key>CCC_BOARD_URL</key>
    <string>http://127.0.0.1:7775</string>
    <key>CCC_CHAT_NO_OPEN</key>
    <string>1</string>
    <key>PATH</key>
    <string>/Users/apple/.local/bin:/Users/apple/.npm-global/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/ccc-chat-server.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/ccc-chat-server.err</string>
</dict>
</plist>
